# files/docker.service.d/10-dockeropts.conf
[Service]
Environment="DOCKER_OPTS=--log-opt max-size=10m --log-opt max-file=1"
