# units/env-lock.path
[Unit]
Description=Bastion - PathLock

[Path]
PathExists=/etc/opsee/bastion-env.sh
